# docker/supervisord.conf
[program:kinesis-local-mock]
command = /app/prepare-env.sh
startsecs = 0
autorestart = false
startretries = 1

[supervisord]
nodaemon=true